# home/z7z8th/.config/autostart/change-ati-clock.desktop
[Desktop Entry]
Type=Application
Name=Change ATI Clock
Exec=/home/z7z8th/gitrepos/z7z8th/scripts/change-ati-clock.sh
#Hidden=True
