# home/z7z8th/.config/autostart/skype.desktop
[Desktop Entry]
Name=skype
Type=Application
Exec=/usr/bin/skype
