# home/z7z8th/.config/autostart/thunderbird.desktop
[Desktop Entry]
Exec=/usr/local/bin/thunderbird
Type=Application
Name=thunderbird
