# home/z7z8th/.config/autostart/xscreensaver.desktop
[Desktop Entry]
Type=Application
Name=xscreensaver
Exec=/usr/bin/xscreensaver
Comment=extensible screen saver framework, plus locking
